# recipes-core/dbus-session/files/dbus-session.socket
[Unit]
Description=D-Bus User Message Bus Socket

[Socket]
ListenStream=/var/run/dbus/session_bus_socket
#ExecStartPost=-/bin/systemctl --user set-environment DBUS_SESSION_BUS_ADDRESS=unix:path=/var/run/dbus/session_bus_socket

#[Install]
#WantedBy=sockets.target
